# systemd/content-generation.service
[Unit]
Description=Generate PULSE posts + SIGNAL videos for all users
After=network.target personal-intelligence-api.service
Wants=personal-intelligence-api.service

[Service]
Type=oneshot
EnvironmentFile=/home/ec2-user/personal-intelligence/.env
# The API (port 8000) serves these jobs and authenticates the same-host call
# with INTERNAL_SECRET. Each job is async (HTTP 202) and fans out over every
# onboarded user via run_for_all_users. Stages run in order: topics first (so
# the pool reflects new memories), then posts (default 3/user), then videos
# (default 1/user). The sleeps let an earlier stage's background work progress.
# The trailing `; echo` makes each step exit 0 so one stage's failure never
# blocks a later stage (posts/videos work off the existing topic pool); the
# logged HTTP status (topics=NNN) is the signal if a stage errored.
ExecStart=/bin/sh -c 'curl -sS -o /dev/null -w "topics=%%{http_code} " -X POST -H "X-Internal-Secret: $INTERNAL_SECRET" http://127.0.0.1:8000/api/content/jobs/topics/run; echo'
ExecStart=/bin/sh -c 'sleep 120'
ExecStart=/bin/sh -c 'curl -sS -o /dev/null -w "posts=%%{http_code} " -X POST -H "X-Internal-Secret: $INTERNAL_SECRET" http://127.0.0.1:8000/api/content/jobs/posts/run; echo'
ExecStart=/bin/sh -c 'sleep 30'
ExecStart=/bin/sh -c 'curl -sS -o /dev/null -w "videos=%%{http_code} " -X POST -H "X-Internal-Secret: $INTERNAL_SECRET" http://127.0.0.1:8000/api/content/jobs/videos/run; echo'
