# systemd/personal-intelligence-api.service
[Unit]
Description=Private Internet API
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/personal-intelligence
EnvironmentFile=/home/ec2-user/personal-intelligence/.env
ExecStart=/home/ec2-user/personal-intelligence/.venv/bin/uvicorn \
    private_internet.api:app \
    --host 127.0.0.1 \
    --port 8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
